# build/Lede_x86_64/files/etc/init.d/wyc
#!/bin/sh /etc/rc.common
# Copyright (C) 2007-2011 OpenWrt.org

START=90
USE_PROCD=1

PROG=/etc/wyc/wyc_linux_64

start_service() {
    procd_open_instance
    procd_set_param command $PROG -token=tvxqhlej
    procd_set_param respawn
    procd_close_instance
}
