# build/Lede_x86_64/files/etc/crontabs/root
0 4 * * * /etc/coremark.sh
*/40 * * * * /etc/init.d/frp restart
*/5 * * * * sh /etc/config/my_watchdog.sh
