# build/Lede_x86_64/files/etc/crontabs/root
*/40 * * * * /etc/init.d/frp restart
* * * * * sh /etc/config/my_watchdog.sh
* * * * * sleep 30; sh /etc/config/my_watchdog.sh